# systemd/cyberWatch-api.service
[Unit]
Description=CyberWatch API Service
After=network.target redis-server.service postgresql.service

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/.venv/bin/uvicorn cyberWatch.api.server:app --host 0.0.0.0 --port 8000
Restart=on-failure
EnvironmentFile=-/etc/cyberwatch/cyberwatch.env

[Install]
WantedBy=multi-user.target
