# systemd/cyberWatch-api.service
[Unit]
Description=CyberWatch API Service
After=network.target redis-server.service postgresql.service

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/.venv/bin/uvicorn cyberWatch.api.server:app --host 0.0.0.0 --port 8000
Restart=on-failure
Environment=CYBERWATCH_PG_DSN=postgresql://postgres:postgres@localhost:5432/cyberWatch
Environment=CYBERWATCH_REDIS_URL=redis://localhost:6379/0
Environment=NEO4J_URI=bolt://localhost:7687
Environment=NEO4J_USER=neo4j
Environment=NEO4J_PASSWORD=neo4j

[Install]
WantedBy=multi-user.target
