# systemd/cyberWatch-worker@.service
[Unit]
Description=CyberWatch Measurement Worker (Instance %i)
After=network.target redis-server.service postgresql.service
Wants=redis-server.service postgresql.service

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@

# Environment variables
Environment="WORKER_ID=%i"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/cyberwatch/cyberwatch.env

# Execution
ExecStart=@ROOT_DIR@/.venv/bin/python -m cyberWatch.workers.worker

# Restart policy
Restart=always
RestartSec=5s
StartLimitInterval=300s
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cyberwatch-worker-%i



[Install]
WantedBy=multi-user.target
