# systemd/cyberWatch-enrichment.service
[Unit]
Description=CyberWatch Enrichment Service
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/.venv/bin/python -m cyberWatch.enrichment.run_enrichment
Restart=on-failure
EnvironmentFile=-/etc/cyberwatch/cyberwatch.env

[Install]
WantedBy=multi-user.target
