# systemd/cyberWatch-worker@.service
[Unit]
Description=CyberWatch Measurement Worker (Instance %i)
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=simple
User=cyberwatch
Group=cyberwatch
WorkingDirectory=/opt/cyberwatch

# Environment variables
Environment="CYBERWATCH_PG_DSN=postgresql://cyberwatch:cyberwatch@localhost/cyberwatch"
Environment="CYBERWATCH_REDIS_URL=redis://localhost:6379/0"
Environment="WORKER_ID=%i"
Environment="PYTHONUNBUFFERED=1"

# Execution
ExecStart=/usr/bin/python3 -m cyberWatch.workers.worker

# Restart policy
Restart=always
RestartSec=5s
StartLimitInterval=300s
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cyberwatch-worker-%i

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cyberwatch/logs

[Install]
WantedBy=multi-user.target
